# krakenrf-krakensdr-doa-git/krakensdr_doa.service
[Unit]
Description=KrakenSDR DoA Startup Service

[Service]
Type=oneshot
ExecStart=krakensdr_doa start

[Install]
WantedBy=default.target
