# krakenrf-heimdall-daq-standalone-git/heimdall.service
[Unit]
Description=KrakenSDR Heimdall Startup Service

[Service]
Type=oneshot
ExecStart=heimdall start

[Install]
WantedBy=default.target
